# ircdcron/ircd.cron
0,10,20,30,40,50 * * * *  /scsi2/server1/codemastr/unreal-devel/ircdcron/ircdchk >/dev/null 2>&1
